# docker/supervisor/conf/demo.ini
[program:demo]
command=php /data/httpd/oms/app/taskmgr/sw.php worker:start -d
stdout_logfile=/proc/1/fd/1
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
numprocs=1
user=root
startsecs=3
startretries=100000